# Device_B/fly-flash-auto.service
[Unit]
Description=Fly-Flash Auto Programming Service
After=network.target multi-user.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/fly-flash-auto.sh
Restart=no
User=root
Group=root
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target